# services/polipo/run
#!/bin/bash

umask 000
exec /sbin/setuser abc polipo -c /config/polipo.conf
