# root/etc/services.d/polipo/run
#!/usr/bin/with-contenv bash

umask 000

exec s6-setuidgid abc polipo -c /config/polipo.conf
